# deploy/systemd/lankeeper-agent.service
[Unit]
Description=LANKeeper Privileged Agent
PartOf=lankeeper.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/lankeeper agent
Restart=always
RestartSec=3

User=root
RuntimeDirectory=lankeeper
RuntimeDirectoryMode=0755

ProtectHome=true
PrivateTmp=true
NoNewPrivileges=false
ProtectSystem=full
ReadWritePaths=/etc/default /boot

[Install]
WantedBy=lankeeper.target
